# scripts/crontab
0,20,40 * * * * root source /var/www/venvs/homepage/bin/activate && cd /var/www/vhosts/homepage/app/ && python manage.py update_feeds --settings=settings_local
